# rekordbox.service
[Unit]
Description=Rekordbox USB Export Server
After=network.target
After=mnt-ssd.mount

[Service]
Type=simple
User=rekordbox
Group=rekordbox

# Paths - adjust to your setup
ExecStart=/usr/local/bin/rekordbox-server \
    --music-dir /mnt/ssd/pre-export \
    --cache-dir /mnt/ssd/rekordbox-cache \
    --socket /run/rekordbox/rekordbox.sock

# Socket directory
RuntimeDirectory=rekordbox
RuntimeDirectoryMode=0755

# Memory limits for Dell Wyse (8GB total)
MemoryMax=2G
MemoryHigh=1G

# Restart on failure
Restart=on-failure
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=rekordbox

[Install]
WantedBy=multi-user.target
